# docker-utils/supervisord/php-fpm7.1.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.1 -F
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true

