# docker-utils/supervisord/rsyslog.conf
[program:rsyslog]
command=service rsyslog restart
